# .config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Restart=alwats
ExecStart=sudo xremap /home/druhan/.config/xremap/config.yml --watch
